# cron/letterboxd.cron
# Weekly refresh every Monday at 02:00 UTC
0 2 * * 1 /usr/local/bin/update_letterboxd.sh
